# hub-components/go-test/sv/mysql/run
#!/bin/bash

/mysql-entrypoint.sh mysqld

